# rootfs/etc/s6-overlay/s6-rc.d/microsoft/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start Microsoft service
# ==============================================================================
flags=()
if bashio::config.true 'update_languages'; then
    flags+=('--update-languages')
fi

if bashio::config.true 'debug_logging'; then
    flags+=('--debug')
fi

exec python3 -m wyoming-microsoft-stt \
    --uri "tcp://0.0.0.0:10300" \
    --subscription-key "$(bashio::config 'subscription_key')" \
    --service-region "$(bashio::config 'service_region')" \
    --download-dir /data \
    ${flags[@]}
